# daily-music/com.xuejin.dailymusic.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.xuejin.dailymusic</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/xuejin/Documents/claude/自动化/daily-music/daily_music.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/xuejin/Documents/claude/自动化/daily-music/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/xuejin/Documents/claude/自动化/daily-music/logs/launchd_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>zh_CN.UTF-8</string>
        <key>HOME</key>
        <string>/Users/xuejin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
