# tuige/com.xuejin.tuige.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.xuejin.tuige</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/xuejin/Documents/claude/自动化/tuige/tuige_main.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/xuejin/Documents/claude/自动化/tuige/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/xuejin/Documents/claude/自动化/tuige/logs/launchd_stderr.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
